# cron.d/update_bbox
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
15 * * * * pirate sh ${BBOX_PATH}/update.sh

